# home/dot_config/systemd/user/hypridle.service.d/custom-config.conf
[Service]
Environment=HLDE_HYPRIDLE_CUSTOM=1
